# nas-web.service
[Unit]
Description=nas-web — NAS Panel Web Server
Documentation=https://github.com/nas-panel/nas-web
After=network.target
Wants=network.target

[Service]
Type=simple

# ── Least-privilege execution ──────────────────────────────────────────────
# Run as a dedicated system account instead of root.
# Create it once after install:
#   useradd --system --no-create-home --shell /usr/sbin/nologin nas-web
User=nas-web
Group=nas-web

# Grant only the capability needed to bind ports < 1024 (80, 443).
# CAP_NET_BIND_SERVICE replaces running as root for port binding.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# ── Process hardening ──────────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# ── Writable paths required by nas-web ────────────────────────────────────
ReadWritePaths=/var/log/nas-web /var/run /etc/nas-web /opt/nas-panel

# ── Execution ─────────────────────────────────────────────────────────────
ExecStart=/usr/local/bin/nas-web /etc/nas-web/nas-web.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nas-web

WorkingDirectory=/etc/nas-web

[Install]
WantedBy=multi-user.target
