# packaging/debian/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop nas-web || true
systemctl disable nas-web || true
