# qnap-like-dashboard/debian/nas-docker-autostart.service
[Unit]
Description=Restore Docker containers on boot
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-autostart.sh start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target