# webserver/debian/postinst
#!/bin/bash

# Create log directory
mkdir -p /var/log/nas-web
chown root:root /var/log/nas-web
chmod 755 /var/log/nas-web

# Reload systemd
systemctl daemon-reload

# Enable and start service
systemctl enable nas-web.service
systemctl start nas-web.service

exit 0
